# railboard.service
[Unit]
Description=railboard OLED train departure board
After=network-online.target
Wants=network-online.target

# This is a TEMPLATE. install.sh fills in @USER@/@RAILBOARD@/@CONF@/@ENV@ and
# writes the result to /etc/systemd/system/railboard.service. To install by hand,
# replace the @...@ placeholders yourself.

[Service]
Type=simple
User=@USER@
# Grant /dev/i2c-* access without running as root (the i2c group owns the device).
SupplementaryGroups=i2c
Environment=PYTHONUNBUFFERED=1
# Put RDM_API_KEY=your_consumer_key in this file. '-' => don't fail if missing.
EnvironmentFile=-@ENV@
ExecStart=@RAILBOARD@ --config @CONF@ --display real
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
